# docker/cron/crontab
# Article Processing Cron Jobs
# ========================================
# Environment variables are passed via Docker environment
# PATH is set to ensure PHP CLI is found

SHELL=/bin/bash
PATH=/usr/local/bin:/usr/local/sbin:/usr/bin:/usr/sbin:/bin:/sbin

# Post-process articles (QA, indexing, mark as indexed)
# Runs every 10 minutes to process articles ingested by hydration worker
*/10 * * * * /app/docker/cron/post_process_articles.sh >> /var/log/cron-post-process.log 2>&1

# Backfill historical articles
# Runs once daily at 2 AM for historical data
0 2 * * * /app/docker/cron/index_articles.sh >> /var/log/cron-backfill.log 2>&1

# Cache latest articles for Redis views
# Runs every 15 minutes to keep cache fresh
*/15 * * * * . /etc/environment && /usr/local/bin/php /app/bin/console app:cache-latest-articles >> /var/log/cron-cache-articles.log 2>&1

# Fetch highlights from Nostr relays to database
# Runs every 30 minutes (reduced from 15 — network-intensive external relay calls)
*/30 * * * * . /etc/environment && /usr/local/bin/php /app/bin/console app:fetch-highlights --limit=200 >> /var/log/cron-fetch-highlights.log 2>&1

# Cache latest highlights for Redis views
# Runs every 30 minutes
*/30 * * * * . /etc/environment && /usr/local/bin/php /app/bin/console app:cache-latest-highlights >> /var/log/cron-cache-highlights.log 2>&1

# Cache media discovery events
# Runs every 6 hours to refresh media content
0 */6 * * * /app/docker/cron/media_discovery.sh >> /var/log/cron-media-discovery.log 2>&1

# Project magazine indices to Magazine entities
# Runs every 30 minutes (reduced from 10 — magazine events are infrequent)
*/30 * * * * /app/docker/cron/project_magazines.sh >> /var/log/cron-magazine-projection.log 2>&1

# Warm Unfold site caches
# Runs every 30 minutes to keep subdomain magazine pages loading quickly.
# Since the graph layer (parsed_reference + current_record), this resolves
# from local DB in milliseconds instead of relay round-trips — kept to
# ensure SWR cache is always pre-warmed so no visitor hits a cold miss.
*/30 * * * * /app/docker/cron/unfold_cache_warm.sh >> /var/log/cron-unfold-cache.log 2>&1

# Graph layer consistency audit
# Runs daily at 3 AM to detect drift between event/article tables and
# graph tables (current_record, parsed_reference). Auto-repairs issues.
0 3 * * * . /etc/environment && /usr/local/bin/php /app/bin/console dn:graph:audit --fix >> /var/log/cron-graph-audit.log 2>&1

# Subscription payment checks (Active Indexing + Updates Pro)
# Checks for zap receipts matching pending invoices; also processes expiry / grace transitions.
*/5 * * * * . /etc/environment && /usr/local/bin/php /app/bin/console active-indexing:check-receipts >> /var/log/cron-zap-receipts.log 2>&1
*/60 * * * * . /etc/environment && /usr/local/bin/php /app/bin/console updates-pro:expire-subscriptions >> /var/log/cron-updates-pro-expire.log 2>&1

